# supervisord.conf
[unix_http_server]
file=/tmp/supervisor_monitor.sock
 
[supervisord]
pidfile=/tmp/supervisord_monitor.pid
logfile_backups=1
 
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
 
[supervisorctl]
serverurl=unix:///tmp/supervisor_monitor.sock
 
[inet_http_server]
port=127.0.0.1:9001
username = favbuy
password = tempfavbuy

[program:gunicorn_django]
directory = /srv/baokuan/
environment = PYTHONPATH=/srv/baokuan/, ENV=PRODUCTION, C_FORCE_ROOT="true"
command = gunicorn_django -w=4 --log-file=/tmp/gunicorn_django.log
autorestart=true